# notifier/software.pep.notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>software.pep.notifier</string>
	<key>ProgramArguments</key>
	<array>
		<string>/var/jb/usr/libexec/pep-native-notifier-launcher</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>ProcessType</key>
	<string>Background</string>
	<key>UserName</key>
	<string>mobile</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/var/jb/usr/bin:/var/jb/bin:/usr/bin:/bin</string>
		<key>PEP_HEADLESS_NOTIFIER</key>
		<string>1</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/var/mobile/Library/Logs/software.pep.notifier.log</string>
	<key>StandardErrorPath</key>
	<string>/var/mobile/Library/Logs/software.pep.notifier.log</string>
</dict>
</plist>
